# rpi-scripts/systemd/backup-rsync.timer
[Unit]
Description=Daily backup at 03:00

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
