# rpi-scripts/systemd/checkweb@.service
[Unit]
Description=Check website %i

[Service]
Type=oneshot
EnvironmentFile=/opt/rpi-ops/.env-rpi-ops
ExecStart=/usr/local/bin/checkweb.sh %i
