# rpi-scripts/systemd/docker-clean.service
[Unit]
Description=Docker prune

[Service]
Type=oneshot
ExecStart=/usr/bin/env bash -c '/usr/local/bin/docker_cleanup.sh'
