# rpi-scripts/systemd/journal-trim.timer
[Unit]
Description=Trim journal every day

[Timer]
OnCalendar=*-*-* 04:10:00
Persistent=true

[Install]
WantedBy=timers.target
