# rpi-scripts/systemd/rpi-health.service
[Unit]
Description=RPi Health summary

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rpi_health.sh
