# rpi-scripts/systemd/smart-scan.timer
[Unit]
Description=Daily SMART scan

[Timer]
OnCalendar=*-*-* 02:40:00
Persistent=true

[Install]
WantedBy=timers.target
